# system_files/etc/systemd/system/ai-desktop-agent.service
[Unit]
Description=Khazar Desktop Agent
After=ai-orchestrator.service
PartOf=khazar.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ai-desktop-agent
Restart=always
RestartSec=5
RuntimeDirectory=khazar
RuntimeDirectoryMode=0750
User=khazar
Group=khazar

[Install]
WantedBy=khazar.target
